# platform/arch/ramdisk/mnt-ramdisk.mount
[Unit]
Description=Mount ramdisk

[Mount]
What=tmpfs
Where=/mnt/ramdisk
Type=tmpfs
Options=size=32G
